# templates/webhook.service.j2
[Unit]
Description=Adnan Webhooks
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=webhook
Restart=on-failure
RestartSec=5
ExecStart=/usr/bin/webhook -verbose -hotreload -hooks /var/webhook/hooks.json

[Install]
WantedBy=multi-user.target
